# modules/backup/vps/timers/backup-vps.timer
[Unit]
Description=Daily VPS backup pull

[Timer]
OnCalendar=*-*-* 16:00:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
